# Debian/PHP/8.1/Fpm-Nginx/programs.conf
[program:nginx]
command=nginx -g 'daemon off;'
autorestart=true
autostart=true

[program:php_fpm]
command=php-fpm8.1 -F
autorestart=true
autostart=true