# Debian/PHP/7.3/Fpm-Nginx/programs.conf
[program:nginx]
command=nginx -g 'daemon off;'
autorestart=true
autostart=true

[program:php_fpm]
command=php-fpm7.3 -F
autorestart=true
autostart=true